# rpisurv.service
[Unit]
Description=Rpisurv Raspberry Pi Surveillance
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/rpisurv
KillMode=control-group
Restart=on-failure
RestartSec=20

[Install]
WantedBy=multi-user.target
